# lucaschmid.net.rebuild.service
[Unit]
Type=simple

[Service]
ExecStart=/var/app/lucaschmid.net/watch-rebuild.sh

